# etc/systemd/system/var-lib-docker.mount
# WARNING: This file was auto generated by container-storage-setup. Do not
# edit it. In the future, this file might be moved to a different location.
[Unit]
Description=Mount dockerlv on /var/lib/docker directory.
Before=docker-storage-setup.service

[Mount]
What=/dev/dockervg/dockerlv
Where=/var/lib/docker
Type=xfs
Options=pquota

[Install]
WantedBy=docker-storage-setup.service
